# res/omni-gpu-agent/omni-gpu-agent.service
[Unit]
Description=omni-gpu-agent
After=network.target
Requles=network.target

[Service]
Type=forking
ExecStart=/opt/omni-gpu-agent/start.sh
ExecStop=/root/omni-gpu-agent/stop.sh
KillMode=none
Restart=on-failure

[Install]
WantedBy=multi-user.target